# ansible/templates/supervisor.conf.j2
; Supervisor configuration for rSpotify Bot
; This file manages the bot as a system service
; Generated by Ansible for {{ environment }} environment

[program:rspotify-bot]
; Command to run the bot
command={{ app_dir }}/venv/bin/python {{ app_dir }}/src/rspotify.py

; Working directory
directory={{ app_dir }}/src

; Run as this user
user={{ app_user }}

; Auto-start on boot
autostart=true

; Auto-restart if crashed
autorestart=true

; Maximum number of restart attempts
startretries=3

; Time to wait before considering the program started
startsecs=10

; SIGTERM for graceful shutdown
stopsignal=TERM

; Time to wait for graceful shutdown before SIGKILL
stopwaitsecs=10

; Redirect stderr to stdout
redirect_stderr=true

; Log files
stdout_logfile={{ app_dir }}/logs/bot_output.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile={{ app_dir }}/logs/bot_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5

; Environment variables
environment=PATH="{{ app_dir }}/venv/bin:/usr/local/bin:/usr/bin:/bin"

; Priority (lower starts first)
priority=999

; Send events for monitoring
events=true